# sump_monitor.service
[Unit]
Description=Sump Pump Water Level Monitor
After=network.target

[Service]
Type=simple
User=root
Environment="SUMP_NOTIFY_SCRIPT=/usr/local/bin/sump_notify.sh"
ExecStart=/usr/local/bin/sump_monitord
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
