# sump_monitor.service
[Unit]
Description=Sump Pump Water Level Monitor
After=network.target

[Service]
Type=simple
User=root
Environment="SUMP_NOTIFY_SCRIPT=/usr/local/bin/sump_notify.sh"
Environment="SUMP_POLL_INTERVAL=1"
Environment="SUMP_WET_THRESHOLD=300"
Environment="SUMP_DRY_THRESHOLD=50"
Environment="SUMP_ADC_SAMPLES=5"
Environment="SUMP_ALARM_DELAY=300"
Environment="SUMP_ALERT_HOLDDOWN=300"
Environment="SUMP_STATE_THRESHOLD=3"
StandardOutput=journal
StandardError=journal
ExecStart=/usr/local/bin/sump_monitord
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
